# docker/supervisor.d/app.ini
[supervisord]
nodaemon=true
user=root

[program:phalmon]
command=php /var/www/html/app/cli.php queue
directory=/var/www/html
user=root
process_name=root
numprocs=1
autostart=true
autorestart=true
redirect_stderr=True
